# supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=/usr/local/bin/wrapdocker

[program:java-jenkins-slave]
command=/bin/bash -c "eval $JENKINS_COMMAND"